# scripts/cookie-refresher.service
[Unit]
Description=YouTube Cookie Refresher Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/your/bot
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Run the scheduler
ExecStart=/usr/bin/python3 /path/to/your/bot/scripts/cookie_scheduler.py --headless --cookies-file=/path/to/your/bot/cookies.txt --profile-dir=/path/to/your/bot/browser_profile --min-interval=1 --max-interval=5

# Restart on failure
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/cookie_refresher.log
StandardError=append:/var/log/cookie_refresher_error.log

[Install]
WantedBy=multi-user.target
